# packages/dipstick/debian/dipstick.service
[Unit]
Description=Dipstick

[Service]
ExecStart=/usr/bin/dipstick
Restart=always

[Install]
WantedBy=multi-user.target
