# packages/dipstick/debian/dipstick.service
[Unit]
Description=Dipstick

[Service]
ExecStart=/usr/bin/dipstick

[Install]
WantedBy=multi-user.target
